# rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bash

source /etc/postgresql/common



exec s6-setuidgid -u ${PG_UID} -g ${PG_GID} /usr/bin/postgres -D ${PG_DATADIR};