# rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid ${PG_UID}:${PG_GID} postgres -- -D ${PG_DATADIR};