# rootfs/etc/s6-overlay/s6-rc.d/live-postgres-config/up
/etc/s6-overlay/s6-rc.d/live-postgres-config/run