# rootfs/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -e

# shellcheck source=../../../conf.d/include/cron
source /etc/conf.d/include/cron

# Start cron
exec /usr/sbin/crond -f -l 4 -L "${CRON_LOG_FILE}" -c "${CRON_CRONTAB_DIR}"
